# cron/crontab
0 1,5,9,13,17,21 * * * /usr/local/bin/generate-sitemap >> /tmp/log 2>&1
0 5 */7 * * /usr/local/bin/localisation-update >> /tmp/log 2>&1
